# switch/switch.service
[Unit]
Description=GPIO poweroff button service

[Service]
ExecStart=/home/osmc/switch/switch.sh

[Install]
WantedBy=multi-user.target
